# rsstorrent-openrc
#!/sbin/runscript
# Copyleft Peter Asplund 2012
# Distributed under the terms of the GNU General Public License v2
# description: rsstorrent
#

# Get function from functions library
source /etc/init.d/functions.sh

depend() {
    need localmount net
}

# Setup variables
NAME=rsstorrent
DAEMON=/usr/bin/${NAME}
PIDFILE=/var/run/${NAME}/${NAME}.pid
LOGFILE=/var/log/${NAME}/${NAME}.log

# Check for an invalid user or one without a home directory
eval USERHOME=~$RSSTORRENT_USER
if [ "${USERHOME#/}" = "${USERHOME}" ]; then
    eerror "${NAME}: The user '$RSSTORRENT_USER' specified in $DEFAULT is invalid."
    exit 1
fi

# Check if group is not specified and assign a proper group
if [ -z $GROUP ]; then
   GROUP="users"
fi

# Check config file
checkconfig() {
	local conf="${USERHOME}/.${NAME}/${NAME}.conf"

	if [ -e "${conf}" ] ; then
		if [ -n "$(find ${USERHOME}/.rsstorrent -maxdepth 1 -name ${NAME}.conf -perm +0004)" ] ; then
			eerror "${conf} must not be world-readable. Run e.g.:"
			eerror "   chmod 640 ${conf}"
			eerror "   chown root:rsstorrent ${conf}"
			return 1
		fi
	else
		eerror "Running rsstorrent as $RSSTORRENT_USER"
		eerror "${USERHOME}/.rsstorrent/rsstorrent.conf is needed to run ${NAME}"
		eerror "There is a sample file included in the package"
		return 1
	fi
}

# Start the service
start() {
        checkconfig || return 1
		local piddir=$(dirname ${PIDFILE})
		if [ ! -d ${piddir} ] ; then
			ebegin "Making ${piddir}"
			mkdir -p ${piddir}
			eend $?
			ebegin "Changing permissions of ${piddir}"
			chown ${RSSTORRENT_USER}:${GROUP} ${piddir}
			eend $?
		fi

		local logdir=$(dirname ${LOGFILE})
		if [ ! -d ${logdir} ] ; then
			ebegin "Making ${logdir}"
			mkdir -p ${logdir}
			eend $?
			ebegin "Changing permissions of ${logdir} to ${RSSTORRENT_USER}:${GROUP}"
			chown ${RSSTORRENT_USER}:${GROUP} ${logdir}
			eend $?
		fi
        ebegin "Starting ${SVCNAME} service: "

        ### Create the lock file ###
        start-stop-daemon --background --start --user ${RSSTORRENT_USER} --name ${NAME} \
			--make-pidfile --pidfile ${PIDFILE} --exec ${DAEMON} -- --logfile ${LOGFILE} -d
        eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --signal USR1 --pidfile ${PIDFILE}
	eend $?
}
